# services/smf-manifest/rsync-loop.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!-- This file is an SMF manifest for 'rsync-loop.sh' script, part of COSas package -->
<service_bundle type='manifest' name='export'>
  <service name='application/rsync-loop' type='service' version='0'>
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    <dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
    <instance name='default' enabled='false'>
      <method_context project=':default' resource_pool=':default' working_directory=':default'>
        <method_credential limit_privileges=':default'
         privileges='basic,net_privaddr' supp_groups=':default'
         group=':default' user='root'/>
        <!-- NOTE: this privilege set allows an unprivileged user account
             to utilize privileged ports like HTTP = 80/tcp -->
      </method_context>
      <exec_method name='start' type='method' exec='/opt/COSas/bin/rsync-loop.sh %m' timeout_seconds='60'>
        <method_context/>
      </exec_method>
      <exec_method name='stop' type='method' exec=':kill' timeout_seconds='30'>
        <method_context/>
      </exec_method>
      <property_group name="startd" type="framework">
            <propval name="duration" type="astring" value="contract"/>
            <propval name="ignore_error" type="astring" value="core,signal"/>
      </property_group>
    </instance>
    <template>
      <common_name>
        <loctext xml:lang='C'>rsync-loop</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
